# configs/cmdd.service.in
[Unit]
Description=@DESCRIPTION@
After=connman.service

[Service]
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/@PROJECT_NAME@

[Install]
WantedBy=multi-user.target
